# deploy/systemd/shinzo-defradb.service
[Unit]
Description=Shinzo Network DefraDB Node
Documentation=https://github.com/sourcenetwork/defradb
After=network.target
Wants=network.target

[Service]
Type=simple
User=shinzo
Group=shinzo
WorkingDirectory=/opt/shinzo
ExecStart=/opt/shinzo/bin/defradb start --store-path=/opt/shinzo/data --url=0.0.0.0:9181
ExecReload=/bin/kill -HUP $MAINPID
KillMode=process
Restart=on-failure
RestartSec=5
TimeoutStopSec=30

# Environment file
EnvironmentFile=/opt/shinzo/.env

# Security settings
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/opt/shinzo/data /opt/shinzo/logs
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=shinzo-defradb

[Install]
WantedBy=multi-user.target
